# SVMC-GAME-2/etc/profile.d/java.sh
pathmunge /usr/java/jdk-17.0.1+12-jre/bin